# stowed/.config/systemd/user/cameractrls-usb-CN09357G8LG009ARAC2KA01_Integrated_Webcam_HD-video-index0.path
[Unit]
Description=Watch /dev/v4l/by-id/usb-CN09357G8LG009ARAC2KA01_Integrated_Webcam_HD-video-index0 and restore controls

[Path]
PathExists=/dev/v4l/by-id/usb-CN09357G8LG009ARAC2KA01_Integrated_Webcam_HD-video-index0

[Install]
WantedBy=paths.target
